# ansible/roles/update-server/files/snapshot-flush.timer
[Unit]
Description=run flush service at midnight utc

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
