# ansible/roles/update-server/templates/worker-restart.service
[Unit]
Description=restart update-worker

[Service]
ExecStart=/bin/systemctl restart update-worker
